# extra/autostart/pasteblock.desktop
[Desktop Entry]
Name=xmousepasteblock
Type=Application
Exec=xmousepasteblock