# extra/autostart/volctl.desktop
[Desktop Entry]
Name=Volctl
Type=Application
Exec=volctl